# supervisor.conf
[supervisord]
nodaemon=true

[program:pmcd]
command=/etc/init.d/pmcd start

[program:pmproxy]
command=/etc/init.d/pmproxy start
